# scripts/vagrant/files/django_supervisor.conf
[program:django]
command=/usr/local/bin/uwsgi --ini /home/vagrant/htdocs/wsgi/wsgi.ini
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True
